# libs/server/scripts/deploy/etc/cron.d/clean_and_update
0 3 * * * root /opt/monorepo/server/clean_and_update.sh
