# libs/server/scripts/deploy/etc/cron.d/docker_vars
@reboot sysctl -w vm.max_map_count=262144
@reboot sysctl -w vm.overcommit_memory=1
