# libs/server/scripts/etc/cron.d/clean_and_update
0 3 * * * root /home/jkonratt/server/clean_and_update.sh
