# libs/server/scripts/etc/cron.d/copy_backups
 0  8 * * * root /home/jkonratt/server/copy_backups.sh
30 20 * * * root /home/jkonratt/server/copy_backups.sh
