# libs/server/scripts/etc/cron.d/copy_backups
 0  8 * * * root rsync -avuH --delete /home/jkonratt/server/backup-service/ /justSSD/backups
30 20 * * * root rsync -avuH --delete /home/jkonratt/server/backup-service/ /justSSD/backups
