# libs/server/scripts/etc/cron.d/create_playlists
0 3 * * * jkonratt /home/jkonratt/server/create_playlists.sh /home/jkonratt/music
