# libs/server/scripts/etc/cron.d/energy_saver_cron
0 7 * * * root /home/jkonratt/server/energy_saver_off.sh
30 22 * * * root /home/jkonratt/server/energy_saver_on.sh
