# libs/server/scripts/etc/cron.d/renovate
0 9 * * * root docker compose -f /home/jkonratt/server/renovate/docker-compose.yml --env-file /home/jkonratt/.env --env-file /home/jkonratt/server/.env --env-file /home/jkonratt/server/renovate/.env up
