# libs/server/scripts/cron.d/energy_saver_cron
30 6 * * * root /home/jkonratt/server/energy_saver_off.sh > /var/log/energy_saver_cron.log 2> /var/log/energy_saver_cron_err.log
0 21 * * * root /home/jkonratt/server/energy_saver_on.sh > /var/log/energy_saver_cron.log 2> /var/log/energy_saver_cron_err.log
